# conf/systemd.service
[Unit]
Description=__APP__
After=network.target

[Service]
WorkingDirectory=__INSTALL_DIR__
EnvironmentFile=__INSTALL_DIR__/.env
Environment=HOST=127.0.0.1
Environment=PORT=__PORT__
User=__APP__
Group=www-data
ExecStart=__NODEJS_DIR__/node .output/server/index.mjs
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
